# .cron
0,30 * * * * /home/fedora41/scripts/git-autopush.sh
1,31 * * * * /home/fedora41/scripts/dotfile-backups.sh
