# .cron
0,30 * * * * /home/jbm/scripts/git-autopush.sh
