# .cron
@reboot /home/jbm/scripts/webstore.pl
0,30 * * * * /home/jbm/scripts/git-autopush.sh
0 */20 * * * /home/jbm/scripts/webstore.pl
