# .cron
0,30 * * * * /home/unix/scripts/git-autopush.sh
1,31 * * * * /home/unix/scripts/dotfile-backups.sh
